# {{ cookiecutter.project_name }}/deploy/supervisor/supervisord.conf
[unix_http_server]
file = /run/supervisor.sock
chown = root:root
chmod = 0700
username = username
password = passwd

[supervisord]
nodaemon = true
nocleanup = true
logfile  =/var/log/supervisord.log
loglevel = warn
childlogdir  =/var/log
user = root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///run/supervisor.sock
username = username
password = passwd

[program:migrate]
command = python manage.py migrate --noinput
directory=/var/www/{{ cookiecutter.project_name }}
autorestart = false
priority = 10
stopwaitsecs = 60
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

[program:nginx]
command = nginx -g "daemon off;"
autostart = true
autorestart = true
priority = 200
stopwaitsecs = 60
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0

{% if SERVER_TYPE == 'uwsgi' -%}

[program:uwsgi]
; command=uwsgi --ini /var/www/uwsgi.ini
command=uwsgi --ini /var/www/uwsgi.ini
user = www-data
directory = /var/www
priority = 400
startsecs = 10
stopwaitsecs = 60
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

{% else -%}

[fcgi-program:asgi]
command = daphne -u /run/daphne/daphne%(process_num)d.sock --fd 0 --access-log - --proxy-headers {{ cookiecutter.project_name }}.asgi:application
process_name=asgi%(process_num)d
socket=tcp://localhost:8000
directory = /var/www/{{ cookiecutter.project_name }}
priority = 400
numprocs = 4
startsecs = 10
stopwaitsecs = 60
user = www-data
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true

{% endif -%}

# [program:dramatiq]
# command = python3 manage.py rundramatiq
# directory = /var/www/{{ cookiecutter.project_name }}
# autostart = true
# autorestart = true
# startsecs = 10
# stopwaitsecs = 60
# stdout_logfile = /dev/stdout
# stdout_logfile_maxbytes = 0
# stderr_logfile = /dev/stderr
# stderr_logfile_maxbytes = 0
